# layout/DEBIAN/postinst
#!/bin/bash

chgrp staff /Library/LaunchDaemons/com.efrederickson.reachapp.fsdaemon.plist
chown root /Library/LaunchDaemons/com.efrederickson.reachapp.fsdaemon.plist
launchctl load -w /Library/LaunchDaemons/com.efrederickson.reachapp.fsdaemon.plist

echo "The best Multitasker you've never used... Until now."